# systemd-user/openclaw-maintenance.service
[Unit]
Description=OpenClaw Daily Maintenance (update + restart + email report)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=PATH=/home/vin/.npm-global/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=/home/vin/.openclaw/workspace
ExecStart=/usr/bin/env bash /home/vin/.openclaw/workspace/scripts/openclaw-maintenance.sh
