# systemd-user/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (v2026.2.23)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/node /home/vin/.npm-global/lib/node_modules/openclaw/dist/index.js gateway --port 18789
Restart=always
RestartSec=5
KillMode=process
Environment=HOME=/home/vin
Environment=TMPDIR=/tmp
Environment=PATH=/home/vin/.local/bin:/home/vin/.npm-global/bin:/home/vin/bin:/home/vin/.volta/bin:/home/vin/.asdf/shims:/home/vin/.bun/bin:/home/vin/.nvm/current/bin:/home/vin/.fnm/current/bin:/home/vin/.local/share/pnpm:/usr/local/bin:/usr/bin:/bin
Environment=OPENCLAW_GATEWAY_PORT=18789
Environment=OPENCLAW_GATEWAY_TOKEN=[REDACTED_SECRET]
Environment=OPENCLAW_SYSTEMD_UNIT=openclaw-gateway.service
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=OPENCLAW_SERVICE_VERSION=2026.2.23

[Install]
WantedBy=default.target
