# systemd-user/openclaw-github-backup.service.d/override.conf
[Service]
Environment="BACKUP_REPO_URL=git@github.com:RBSDigital/OpenclawBackup.git"
Environment="BACKUP_REPO_BRANCH=main"
Environment="GOG_BIN=/home/vin/.local/bin/gog"
Environment="GOG_ACCOUNT=vrbs940054@gmail.com"
Environment="REPORT_EMAIL=vrbs940054@gmail.com"
Environment="PATH=/home/vin/.local/bin:/home/vin/.npm-global/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
